# image_downloader/image_downloader.service
[Unit]
Description=Image Downloader Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/image_downloader
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/image_downloader/image_downloader.py --urls /home/YOUR_USERNAME/image_downloader/urls.txt --output /home/YOUR_USERNAME/image_downloader/images --interval 30 --workers 5 --timeout 5
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/YOUR_USERNAME/image_downloader/logs/service.log
StandardError=append:/home/YOUR_USERNAME/image_downloader/logs/service.log

[Install]
WantedBy=multi-user.target